# scripts/local/launchd/xyz.alphaxbot.mineru-api-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  MinerU pipeline 专用 FastAPI 服务 LaunchDaemon（端口 7862，.venv-pipeline）。
  与 mineru-api（vlm，7861）并存：pipeline 依赖 transformers 4.x，与 mlx 路径的
  transformers 5.x 冲突，必须独立 venv + 独立服务，详见 perf 文档「双服务架构」。
  安装位置：/Library/LaunchDaemons/xyz.alphaxbot.mineru-api-pipeline.plist
  日志：/var/log/mineru-api-pipeline.{out,err}.log

  ⚠️ TCC 完全磁盘访问：本服务的 venv python 路径是
     /Users/macmini/Documents/minerU/.venv-pipeline/bin/python3.12
     与主服务的 .venv python 是【不同二进制】，必须【单独】再加一次 FDA，
     否则 launchd（root）子进程读 ~/Documents 受保护目录会 PermissionError。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>xyz.alphaxbot.mineru-api-pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/macmini/Documents/minerU/scripts/local/mineru-api-pipeline.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>UserName</key>
  <string>macmini</string>
  <key>GroupName</key>
  <string>staff</string>

  <key>WorkingDirectory</key>
  <string>/Users/macmini/Documents/minerU</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/macmini</string>
    <key>PATH</key>
    <string>/Users/macmini/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/var/log/mineru-api-pipeline.out.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/mineru-api-pipeline.err.log</string>

  <!-- 崩溃后 5 秒再重启，避免疯狂打转 -->
  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
